# .config/systemd/user/pulseUI-wifi.service
[Unit]
BindsTo=wifi.target

[Service]
Type=simple
Restart=always
ExecStart=/opt/pulsesecure/bin/pulseUI

[Install]
WantedBy=wifi.target
